# com.zrnet.backupcovers.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zrnet.backupcovers</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>/Users/chsa911/zrnet-clean/backup-covers-incremental.sh</string>
    </array>

    <!-- Runs once a day at 03:00 local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Don't fire immediately on load; wait for the scheduled time -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/chsa911/p_backup/launchd-backup.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/chsa911/p_backup/launchd-backup.log</string>
</dict>
</plist>
